# config/supervisor.http-router.conf
[program:http-router]
priority=10
directory=/usr/local/lib/http-router
command=bundle exec unicorn
autostart=true
autorestart=true
stopsignal=QUIT
