# config/supervisor.nginx.conf
[program:nginx]
priority=20
directory=/tmp
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
user=root
autostart=true
autorestart=true
